# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:h2o]
command=/usr/local/h2o/bin/h2o -m daemon -c /etc/h2o/h2o.conf
autorestart=true
